# .systemd/opencode-notification-worker.service
[Unit]
Description=OpenCode Notification Worker
After=network.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=root
WorkingDirectory=/root/projekt/akm-bridge
ExecStart=/root/.bun/bin/bun run src/notifications/worker.ts
EnvironmentFile=-/root/projekt/akm-bridge/.env.notifications
Restart=on-failure
RestartSec=10
NoNewPrivileges=true
PrivateTmp=true
UMask=0077

[Install]
WantedBy=multi-user.target
