# .systemd/opencode-sli-collector.timer
[Unit]
Description=OpenCode SLI Collector — collect metrics every 60 seconds
Documentation=https://github.com/przemyslawzielinski87-debug/akm-bridge

[Timer]
OnBootSec=60
OnUnitActiveSec=60
AccuracySec=5
RandomizedDelaySec=5
Persistent=false

[Install]
WantedBy=timers.target