# .systemd/opencode-update-check.timer
[Unit]
Description=Weekly OpenCode Update Check

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
